# data/indicator-power.service.in
[Unit]
Description=Indicator Power
PartOf=graphical-session.target
PartOf=indicators.target

[Service]
ExecStart=@pkglibexecdir@/indicator-power-service
Restart=on-failure

[Install]
WantedBy=indicators.target
