# data/indicator-power.service.in
[Unit]
Description=Indicator Power
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/indicator-power-service
Restart=on-failure
